# dockercron/crontab
*/30 * * * * root python /app/mysql_to_csv.py > /proc/1/fd/1 2>/proc/1/fd/2
